# prometheus-xentop.service.in
[Unit]
Description=Export Xen statistics to Prometheus

[Service]
EnvironmentFile=-@SYSCONFDIR@/default/@NAME@
ExecStart=@BINDIR@/@NAME@ $ARGS

[Install]
WantedBy=multi-user.target
